# projects/prometheus/files/prometheus.conf
[program: prometheus]
directory=/var/praekelt/prometheus/
command =/var/praekelt/prometheus/prometheus 
user = root

